# lightdm/.local/etc/profile.d/20-default-applications.sh
# Default applications
# ====================

export BROWSER="firefox"
export EDITOR="vim"
export VISUAL="vim"

export GIT_EDITOR=$VISUAL

export PAGER="less"
export MANPAGER="less -X"
export LESSCHARSET="utf-8"
