# dotfiles/.config/profile.d/99-startx.sh
# ~/.config/profile.d/99-startx.sh
# ================================

# [ -z "$DISPLAY" ] && [ "$XDG_VTNR" == "1" ] && startx
